# supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn app:app --bind 0.0.0.0:8000
autorestart=true

[program:bot]
command=python bot.py
autorestart=true
